# tests/containers/workbenches/jupyterlab/mysql/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/usr/sbin/mysqld
autostart=true
autorestart=true
user=mysql

[program:saslauthd]
command=/usr/sbin/saslauthd -a pam
autostart=true
autorestart=true
